# crontab
# generate/deliver Nightly at 9:59pm CDT
59 21 * * * rake generate issue:deliver
